# container-monitor.service
[Unit]
Description=Credit Checker Container Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/creditcheckerlive
ExecStart=/usr/bin/python3 /root/creditcheckerlive/monitor-containers.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
