# creditchecker-flask.service
[Unit]
Description=Credit Checker Flask Web Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/creditcheckerlive
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/root/creditcheckerlive
ExecStart=/usr/bin/python3 /root/creditcheckerlive/web_server.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=creditchecker-flask

# Security settings (relaxed for root directory access)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
